# root/etc/s6-overlay/s6-rc.d/transcode_stream/run
/etc/s6-overlay/scripts/acquire_slinger_up.sh
