# root/etc/s6-overlay/s6-rc.d/acquire_slinger/up
/etc/s6-overlay/scripts/acquire_airconnec_up.sh
